# deploy/spaans.service
[Unit]
Description=Spaans oefen-webapp
After=network.target

[Service]
Type=simple
User=wouter
WorkingDirectory=/home/wouter/spaans
ExecStart=/home/wouter/spaans/.venv/bin/uvicorn --factory app.main:create_app --host 127.0.0.1 --port 8100
Restart=on-failure
RestartSec=3

# Afscherming
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/wouter/spaans/data

[Install]
WantedBy=multi-user.target
